# files/crontab/storagemgr
#Copy to backup server daily by default
45 2 * * * /app/scripts/copy-backup.sh >> /var/log/cron.log 2>&1
# An empty line is required at the end of this file for a valid cron file.
